# supervisord.conf
[supervisord]
nodaemon=true

[program:webssh]
command=/usr/sbin/sshd -D

[program:nginx]
command=/usr/sbin/nginx